# home/.config/systemd/user/break-timer.service
[Unit]
Description=Break timer

[Service]
Type=simple
ExecStart=/home/mazhar/bin/break-timer
Restart=no

[Install]
WantedBy=default.target
